# Marketplaces - Archive/Solaris DNM/mm-catalog/deploy-config/supervisord.conf
[supervisord]
logfile=/etc/supervisor/logs/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=5MB         ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:laravel-sync]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan rabbitmq:consumer --queue='shop_events' --exchange='catalog' --routingKey='shop_events'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB

[program:laravel-stats]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan rabbitmq:consumer --queue='statistic' --exchange='catalog' --routingKey='statistic'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB
[program:laravel-disputes]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan rabbitmq:consumer --queue='disputes' --exchange='catalog' --routingKey='disputes'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/disputes.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB
[program:laravel-sync-goods]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan rabbitmq:consumer --queue='good_events' --exchange='catalog' --routingKey='good_events'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB
[program:laravel-sync-orders]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan rabbitmq:consumer --queue='order_events' --exchange='catalog' --routingKey='order_events'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stopwaitsecs=3600
stdout_logfile_maxbytes=5MB
